# dist/software_opensuse_org_prometheus_exporter.service
[Unit]
Description=The prometheus exporter for the software.opensuse.org app

[Service]
User=soorun
WorkingDirectory=~
ExecStart=/usr/bin/bundle exec prometheus_exporter

[Install]
WantedBy=multi-user.target

